# supervisor/logstash.conf
[program:logstash]
command=/opt/logstash/bin/logstash agent -f /opt/logstash.conf
stdout_logfile=/dev/stdout.log
redirect_stderr=true
